# zerohid.service
[Unit]
Description=Zerohid
[Service]
ExecStart=/root/zerohid/zerohid.sh
[Install]
WantedBy=multi-user.target
